# build/dbgp-proxy/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:dbgp-proxy]
command=/usr/local/bin/pydbgpproxy -l DEBUG -i %(ENV_DBGP_IDE_PORT)s -d %(ENV_DBGP_DEBUG_LISTEN)s

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
